# root/etc/services.d/lemur/run
#!/usr/bin/with-contenv bash

exec lemur --config=/config/lemur/lemur.conf.py start -b 0.0.0.0:8000