# root/etc/periodic/15min/lemur_sync
#!/usr/bin/with-contenv bash

lemur -c /config/lemur/lemur.conf.py source sync -s all